# 2024/crontab
# m h  dom mon dow   command
0 0 * * Thu /home/cial1/score.sh > /home/cial1/score.log 2>&1
